# systemd-units/update-squid-ads.service
[Unit]
Description=Update Squid ad blocking list (HaGeZi PRO)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-squid-ads.sh

[Install]
WantedBy=multi-user.target
