# systemd-units/squid.service.d/override.conf
[Unit]
# Squid doit démarrer après Caddy (DoH) pour que la résolution DNS soit opérationnelle
After=caddy.service
